# chef-repository/cookbooks/DICETraffic/templates/default/generic.timer.erb
[Unit]
Description=<%= @description %>

[Timer]
OnCalendar=<%= @calendar %>
Unit=<%= @unit %>

[Install]
WantedBy=multi-user.target
