# virtual_machines/pytraffic_vm/systemd/lpp_minutely.timer
[Unit]
Description=Send lpp live arrival to Kafka.

[Timer]
OnCalendar=*-*-* *:*:00
Unit=lpp_minutely.service

[Install]
WantedBy=multi-user.target