# blueprint/chef-repository/cookbooks/DICETraffic/templates/default/generic.service.erb
[Unit]
Description=<%= @description %>

[Service]
Type=simple
ExecStart=/usr/local/bin/pytraffic <%= @parameters %> --config <%= @config_path %>/local.conf

[Install]
WantedBy=multi-user.target
